# apps-conf-files/conky/conky.service
[Unit]
Description=Starts Conky service

[Service]
ExecStart=/home/$USER/.concky/conky-startup.sh

[Install]
WantedBy=multi-user.target
